# airootfs/etc/systemd/system/reflector.service.d/archiso.conf
[Service]
Restart=on-failure
RestartSec=10
